# deploy/systemd/consilium-backup.service
[Unit]
Description=Consilium backup to B2
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/consilium-backup.sh
StandardOutput=journal
StandardError=journal
